# supervisor/dockerd.conf
[program:dockerd]
command=/usr/local/bin/dockerd
autostart=true
startretries=5
startsecs=10
stderr_logfile=/var/log/dockerd.err.log
stdout_logfile=/var/log/dockerd.out.log
