# supervisor/webserver.conf
[program:webserver]
directory=/app
command=/usr/bin/node ./index.js
autostart=true
autorestart=true
stderr_logfile=/var/log/webserver.err.log
stdout_logfile=/var/log/webserver.out.log
